# share/claude-cost-collect.service.in
[Unit]
Description=Collect Claude Code usage data

[Service]
Type=oneshot
ExecStart=/bin/bash %%COLLECT_SCRIPT%%
