# share/com.claude-cost.collect.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-cost.collect</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>%%COLLECT_SCRIPT%%</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>%%HOUR%%</integer>
        <key>Minute</key>
        <integer>%%MINUTE%%</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>%%LOG_DIR%%/collect.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>%%LOG_DIR%%/collect.stderr.log</string>
</dict>
</plist>
